# video-web/doc/supervisor/ssserver.conf
[program:ssserver]
command=/root/.nvm/versions/node/v18.19.0/bin/node ssserver
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/ssserver.stdout.log